# deploy/tavrn.service
[Unit]
Description=tavrn.sh SSH tavern
After=network.target

[Service]
Type=simple
User=tavrn
WorkingDirectory=/home/tavrn/tavrn
EnvironmentFile=-/home/tavrn/tavrn/.env
Environment=TAVRN_PORT=22
ExecStart=/home/tavrn/tavrn/tavrn --web-audio
Restart=always
RestartSec=3
LimitNOFILE=65535
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
